# gnap/filesystem-overlay/common/etc/init.d/fix-permissions
#!/sbin/runscript

start() {
	einfo "Chowning user homedirs"
	for a in ${USERS}
	do
		echo $a
		chown -R ${a}:${a} ~${a}
	done
}
